# .claude/scripts/launchd/com.kriscard.heartbeat.github.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN"
  "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
    <key>Label</key>
    <string>com.kriscard.heartbeat.github</string>
    <key>ProgramArguments</key>
    <array>
        <string>/Users/kriscard/.local/bin/uv</string>
        <string>run</string>
        <string>/Users/kriscard/.dotfiles/.claude/scripts/heartbeat.py</string>
        <string>--mode</string>
        <string>github</string>
    </array>
    <key>StartCalendarInterval</key>
    <array>
        <dict><key>Weekday</key><integer>1</integer><key>Hour</key><integer>9</integer><key>Minute</key><integer>30</integer></dict>
        <dict><key>Weekday</key><integer>2</integer><key>Hour</key><integer>9</integer><key>Minute</key><integer>30</integer></dict>
        <dict><key>Weekday</key><integer>3</integer><key>Hour</key><integer>9</integer><key>Minute</key><integer>30</integer></dict>
        <dict><key>Weekday</key><integer>4</integer><key>Hour</key><integer>9</integer><key>Minute</key><integer>30</integer></dict>
        <dict><key>Weekday</key><integer>5</integer><key>Hour</key><integer>9</integer><key>Minute</key><integer>30</integer></dict>
    </array>
    <key>StandardOutPath</key>
    <string>/tmp/heartbeat-github.log</string>
    <key>StandardErrorPath</key>
    <string>/tmp/heartbeat-github.err</string>
    <key>RunAtLoad</key>
    <false/>
</dict>
</plist>
